# wallpaper.service
[Unit]
Description=Hourly wallpaper

[Service]
Type=wallpaper
ExecStart=dwall -s lakeside
