# wallpaper.service
[Unit]
Description=Hourly wallpaper

[Service]
Type=wallpaper
ExecStart=/home/ren/.wallpaper
